# infra/vps/systemd/alerts-ua-import.service
[Unit]
Description=Alerts UA XLSX import
After=network.target

[Service]
Type=oneshot
User=alerts-ua
Group=alerts-ua
WorkingDirectory=/srv/alerts-ua/app/current
EnvironmentFile=/srv/alerts-ua/env/.env.worker
ExecStart=/srv/alerts-ua/runtime/node/bin/node dist/workers/import.js
