# infra/vps/systemd/alerts-ua-poll.service
[Unit]
Description=Alerts UA polling cycle
After=network.target

[Service]
Type=oneshot
User=alerts-ua
Group=alerts-ua
WorkingDirectory=/srv/alerts-ua/app/current
EnvironmentFile=/srv/alerts-ua/env/.env.worker
ExecStart=/srv/alerts-ua/runtime/node/bin/node dist/workers/poll.js
