# infra/vps/systemd/alerts-ua-telegram-ingest.timer
[Unit]
Description=Run Alerts UA Telegram ingestion every 15 seconds

[Timer]
OnBootSec=45
OnUnitActiveSec=15
Unit=alerts-ua-telegram-ingest.service

[Install]
WantedBy=timers.target
